# master/ansible/examples/fe-be-db/app.service
Description=App Service
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/python /tmp/app.py

[Install]
WantedBy=multi-user.target
